# provisioning-config/etc/cron.d/puppet.cron
# Get config and run puppet every hour. 
@reboot root /usr/bin/runpuppet
3 * * * * root /usr/bin/runpuppet