# etc/cron.d/puppet.cron
# Get config and run puppet every 30 minutes. 
@reboot root /usr/bin/runpuppettwice
*/30 * * * * root /usr/bin/runpuppet
