# etc/cron.d/puppet.cron
# Get config and run puppet every 6 minutes. 
@reboot root /usr/bin/runpuppettwice
*/6 * * * * root /usr/bin/runpuppet
