# etc/cron.d/puppet.cron
# Get config and run puppet every 10 minutes. 
@reboot root /usr/bin/runpuppet
*/10 * * * * root /usr/bin/runpuppet