# etc/cron.d/puppet.cron
# Run puppet every hour and on reboots. 
@reboot /usr/bin/runpuppet
7 * * * * root /usr/bin/runpuppet